# cookbooks-rails/security/templates/default/rc.local
#!/bin/sh -e
#
# rc.local
sh /etc/iptables.sh

# start nginx server
# service passenger start

# start redis server
# service redis-server start

exit 0